# glanus-agent/installers/linux/glanus-agent.service
[Unit]
Description=Glanus Agent
Documentation=https://docs.glanus.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=GLANUS_DAEMON=1
# Best-effort X11 environment discovery for screen capture / input.
# Runs *before* the agent and writes DISPLAY / XAUTHORITY into a runtime
# env file. If no graphical session is active (e.g. server with no X),
# the file is written empty and the agent runs without capture (the
# runtime logs an error in that case so it's diagnosable).
EnvironmentFile=-/run/glanus-agent/x11.env
ExecStartPre=+/usr/bin/glanus-agent-detect-x11
ExecStart=/usr/bin/glanus-agent --daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# RuntimeDirectory + StateDirectory tell systemd to (re)create these
# paths on every start and chown them to the unit's User=. Using
# ReadWritePaths= here previously caused a 226/NAMESPACE failure on
# fresh installs and reboots because /run is tmpfs and /run/glanus-agent
# did not yet exist when systemd was setting up the mount namespace.
# RuntimeDirectoryPreserve=yes keeps /run/glanus-agent across restarts
# of the unit (but not across reboots — tmpfs semantics) so the
# ExecStartPre-written x11.env survives a Restart=always cycle.
RuntimeDirectory=glanus-agent
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
StateDirectory=glanus-agent
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
